# src/research/cron/crontab
* * * * * root bash /bob-research/cron/test_crontab.sh
* * * * * root bash /bob-research/cron/run_heatmap.sh
