# src/research/cron/crontab
# * * * * * /bob-research/cron/test_crontab.sh > /proc/1/fd/1 2>/proc/1/fd/2
0 * * * * /bob-research/cron/run_heatmap.sh > /proc/1/fd/1 2>/proc/1/fd/2
